# states/gcc/files/locales/generate-locales.service
[Unit]
Description=Generate Locales at boot

[Service]
Type=oneshot
ExecStart=/usr/bin/locale-gen

[Install]
WantedBy=mutli-user.target
